# meta-intel-extras/recipes-support/grub/grub-efi/grubenv-copy.service
[Unit]
Description=Grub Environment Copy

[Service]
ExecStart=/usr/bin/grubenv-copy.sh

[Install]
WantedBy=multi-user.target

