# dynamic-layers/acrn-layer/recipes-vms/acrn-vm-init/acrn-vm-init/acrn-vm-init.service
[Unit]
Description=ACRN Initialization procedure

[Service]
Type=oneshot
ExecStart=/usr/share/acrn/scripts/acrn-vm-init.sh

[Install]
WantedBy=multi-user-target
