# recipes-multimedia/pulseaudio/pulseaudio/pulseaudio.service
[Unit]
Description=PulseAudio System Daemon
 
[Install]
WantedBy=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio --system
