# recipes-config/systemd-resource-control/slices/connectivity.slice
[Unit]
Description=Connectivity Slice
Before=slices.target
