# recipes-config/systemd-resource-control/slices/ui.slice
[Unit]
Description=UI Slice
Before=slices.target
